# myapi.service
# File path: /etc/systemd/system/myapi.service
# sudo systemctl daemon-reload
# sudo systemctl start myapi
# sudo systemctl enable myapi
# sudo systemctl status myapi
[Unit]
Description=My API Service
After=network.target

[Service]
User=ubuntu  # 注意修改用户名
WorkingDirectory=/path/to/myapi  # 注意修改路径
ExecStart=/path/to/myapi/myapi --port 3000 --log-dir /path/to/myapi/logs  # 注意修改路径和端口号
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
